# api/configs/deepflow.target
[Unit]
Description=Deepflow

[Install]
WantedBy=multi-user.target
